# root/etc/services.d/mjpg-streamer/run
#!/usr/bin/with-contenv bash

if [ -n "$MJPEG_STREAMER_INPUT" ]; then
  echo "Deprecation warning: the environment variable '\$MJPEG_STREAMER_INPUT' was renamed to '\$MJPG_STREAMER_INPUT'"

  MJPG_STREAMER_INPUT=$MJPEG_STREAMER_INPUT
fi

if ! expr "$MJPG_STREAMER_INPUT" : ".*\.so.*" >/dev/null; then
  MJPG_STREAMER_INPUT="input_uvc.so $MJPG_STREAMER_INPUT"
fi

# Loop through all camera devices and start mjpg_streamer for each
IFS=',' read -ra CAMERA_DEVS <<<"$CAMERA_DEV"
COUNTER=0
for CAMERA in "${CAMERA_DEVS[@]}"; do
  # start from port 8080 and increment for each camera
  PORT=$((8080 + "$COUNTER"))
  ((COUNTER++))
  echo "Starting mjpg_streamer for camera device $CAMERA on port $PORT"

  # Start mjpg_streamer
  mjpg_streamer -b \
    -i "/usr/local/lib/mjpg-streamer/${MJPG_STREAMER_INPUT} -d ${CAMERA}" \
    -o "/usr/local/lib/mjpg-streamer/output_http.so -w /usr/local/share/mjpg-streamer/www -p ${PORT}"
done
