# camera/etc/services.d/haproxy/run
#!/usr/bin/execlineb -P

haproxy -p /tmp/haproxy.pid -f /etc/haproxy/haproxy.cfg -db
