# root/etc/services.d/octoprint/run
#!/usr/bin/with-contenv sh

exec octoprint serve --iknowwhatimdoing --host 0.0.0.0 --port 5000 --basedir /octoprint/octoprint
